# bin/gateway.service
[Unit]
Description=Startup Interfaces and then Gateway server/applications
After=network.target

[Service]
Type=forking
ExecStart=/home/pi/gateway-applications/bin/gateway.start
Restart=always

[Install]
WantedBy=multi-user.target
