# bin/gateway.service
[Unit]
Description=Startup Interfaces and then Gateway server/applications

[Service]
Type=forking
ExecStart=/home/pi/gateway-applications/bin/gateway.start
Restart=on-failure

[Install]
WantedBy=multi-user.target
